# document-api/supervisord.conf
[supervisord]
nodaemon=true
logfile=/server/supervisord.log
childlogdir=/server/

[program:document-api]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
command=/server/document-api
