# docker/supervisord-keyring-api.conf
[supervisord]
nodaemon=true
logfile=/server/supervisord.log
childlogdir=/server/

[program:keyring-api]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
command=/server/keyring-api
